# deploy/systemd/storylens-api.service
[Unit]
Description=Story Lens API (Bun/Elysia)
Documentation=https://github.com/Hussain7Abbas/storylens-backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=storylens
Group=storylens
WorkingDirectory=/opt/storylens/apps/backend
Environment=NODE_ENV=production
EnvironmentFile=/opt/storylens/apps/backend/.env
ExecStart=/usr/local/bin/bun run src/main.ts
Restart=on-failure
RestartSec=5
KillSignal=SIGINT

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
